# bot.service
[Unit]
Description=Bot Service
After=network.target

[Service]
User=root
WorkingDirectory=/root/san/bot/Bot
ExecStart=/usr/bin/python3 main.py
Restart=always

[Install]
WantedBy=multi-user.target
